# tests/golden/webabs/webabs-pull.service
[Unit]
Description=webabs: git pull + rebuild (samosbor)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/samosbor pull webabs
