# tests/golden/go-demo/demo-config.service
[Unit]
Description=demo: restart on config change (samosbor)

[Service]
Type=oneshot
ExecStart=systemctl --user try-restart demo.service
