# tests/golden/go-demo/demo.service
[Unit]
Description=demo (managed by samosbor)

[Service]
ExecStart=/home/user/.local/state/samosbor/demo/current/demo --listen :8080 --verbose
WorkingDirectory=/home/user/.local/state/samosbor/demo
Environment="DEMO_TOKEN=s3cr3t"
Environment="GREETING=hello world"
Environment="TRICKY=50%% \"q\" $dollar \\back"
EnvironmentFile=/etc/demo/env
Restart=always
RestartSec=2
TimeoutStopSec=30

[Install]
WantedBy=default.target
