# tests/golden/legacy/legacy.service
[Unit]
Description=legacy (managed by samosbor)

[Service]
ExecStart=/home/user/.local/bin/legacy --port 9090 --data "/home/user/my data"
WorkingDirectory=/srv/legacy-data
Restart=always
RestartSec=2
TimeoutStopSec=30

[Install]
WantedBy=default.target
