# tests/golden/webapp/webapp-pull.service
[Unit]
Description=webapp: git pull + rebuild (samosbor)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/samosbor pull webapp
